# scripts/distribution/macOS-package/template/payload/Library/Concordium Node/LaunchDaemons/software.concordium.mainnet.node-collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <key>Label</key>
    <string>software.concordium.mainnet.node-collector</string>

    <key>Program</key>
    <string>/usr/local/bin/concordium-node-collector</string>

    <key>EnvironmentVariables</key>
    <dict>

        <!-- Public node name showed on network dashboard. -->
        <key>CONCORDIUM_NODE_COLLECTOR_NODE_NAME</key>
        <string>__NODE_NAME__</string>

        <!-- Dashboard to post collected data to. -->
        <key>CONCORDIUM_NODE_COLLECTOR_URL</key>
        <string>https://dashboard.mainnet.concordium.software/nodes/post</string>

        <!-- Interval in milliseconds to sleep between runs of the collector. -->
        <key>CONCORDIUM_NODE_COLLECTOR_COLLECT_INTERVAL</key>
        <string>5000</string>

        <!-- gRPC host to collect from. -->
        <key>CONCORDIUM_NODE_COLLECTOR_GRPC_HOST</key>
        <string>http://localhost:10000</string>

        <!-- Net name used for logging. Logs can be viewed in Console.app by
             filtering on 'software.concordium.<net-name>.node' where
             '<net-name>' is the value set here. -->
        <key>CONCORDIUM_NODE_COLLECTOR_MACOS_NET_NAME</key>
        <string>mainnet</string>

    </dict>

    <!-- Runs when loaded and is kept alive. -->
    <key>KeepAlive</key>
    <true />

</dict>
</plist>
